# ros.ini

[program:ros]
directory = /data1/ros
command = /data1/ros/ros.sh
stdout_logfile = /dev/null
stderr_logfile = /dev/null
autorestart = true
startretries = 30

